# deploy/supervisor.conf
; ─── OpenSales Supervisor Config ───────────────────────────────
; Joylash: /etc/supervisor/conf.d/opensales.conf
; Reload: sudo supervisorctl reread && sudo supervisorctl update

[program:opensales-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/saleflow/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/opensales-queue.log
stopwaitsecs=3600

[program:opensales-ssr]
process_name=%(program_name)s
command=/var/www/saleflow/deploy/ssr-start.sh
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/opensales-ssr.log
stopasgroup=true
killasgroup=true
environment=HOME="/var/www",NVM_DIR="/var/www/.nvm"

; NOTE: Scheduler cron orqali (system crontab) ishlaydi:
;   * * * * * www-data php /var/www/saleflow/artisan schedule:run >> /dev/null 2>&1
; Cron va supervisor schedule:work — ikkalasi parallel ishlamasligi kerak.
